# repos/system_upgrade/common/actors/initramfs/mount_units_generator/files/bundled_units/boot.mount
[Unit]
DefaultDependencies=no
Before=local-fs.target
After=sysroot-boot.target
Requires=sysroot-boot.target

[Mount]
What=/sysroot/boot
Where=/boot
Type=none
Options=bind
