# repos/system_upgrade/common/actors/initramfs/mount_units_generator/tests/files/unit-nofail.mount
# Automatically generated by systemd-fstab-generator

[Unit]
Documentation=man:fstab(5) man:systemd-fstab-generator(8)
SourcePath=/etc/fstab
Before=local-fs.target
Requires=systemd-fsck@dev-disk-by\x2duuid-1ca63aee\x2d37b5\x2d42b7\x2db1a5\x2df752ba09e010.service
After=systemd-fsck@dev-disk-by\x2duuid-1ca63aee\x2d37b5\x2d42b7\x2db1a5\x2df752ba09e010.service
After=blockdev@dev-disk-by\x2duuid-1ca63aee\x2d37b5\x2d42b7\x2db1a5\x2df752ba09e010.target

[Mount]
What=/dev/disk/by-uuid/1ca63aee-37b5-42b7-b1a5-f752ba09e010
Where=/boot
Type=ext4
Options=nofail,noexec
