# repos/system_upgrade/el7toel8/actors/createresumeservice/files/leapp_resume.service
[Unit]
Description=Temporary Leapp service which resumes execution after reboot
After=default.target
DefaultDependencies=no

[Service]
Type=oneshot
ExecStart=/usr/bin/leapp upgrade --resume
StandardOutput=journal+console
