# repos/system_upgrade/el7toel8/actors/createresumeservice/files/leapp_resume.service
[Unit]
Description=Temporary Leapp service which resumes execution after reboot
After=default.target
DefaultDependencies=no
After=dbus.service
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
# FIXME: this is temporary workround for Python3
ExecStart=/root/tmp_leapp_py3/leapp3 upgrade --resume
StandardOutput=journal+console
